# scripts/chatbox.service
[Unit]
Description=Raspberry Pi KI-Chatbox
After=network-online.target sound.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/mariano/VoiceChatbox
EnvironmentFile=/home/mariano/VoiceChatbox/.env
ExecStart=/home/mariano/VoiceChatbox/.venv/bin/python -m src.main --live-recognition --vosk --chatgpt
Restart=on-failure
RestartSec=2

# Audio permissions
User=mariano
Group=mariano

[Install]
WantedBy=multi-user.target
